# roles/media/files/mopidy.service
[Unit]
Description=Mopidy
After=network.target sound.target

[Service]
ExecStart=mopidy

[Install]
WantedBy=default.target
